# sagemaker/mlflow-sagemaker/user-data/mlflow_supervisord.conf
[program:mlflow]
command=/home/ubuntu/opt/mlflow_server.sh
directory=/home/ubuntu
user=ubuntu
